# kdc-data/etc/systemd/system/getstats@.service
[Unit]
Description=Getstats Per-Connection Server

[Service]
ExecStart=-/usr/local/sbin/getstats.sh
User=root
Group=root
StandardInput=socket
